# cron-scripts/crontab
# look for idp config changes twice per hour
7,37 * * * * /opt/container-scripts/update-configs.sh
