# cron-scripts/crontab
# look for idp config changes twice per hour
7,37 * * * * /opt/cron-scripts/update_configs.sh
